# process-supervisor/s6-overlay/s6-rc.d/node/run
#!/bin/sh
# Next.js standalone server; PORT=3000 HOSTNAME=127.0.0.1 come from the image ENV.
exec node /app/web/server.js
