# process-supervisor/s6-overlay/s6-rc.d/uvicorn/run
#!/bin/sh
# uvicorn binds loopback only; --workers 1 keeps the in-process APScheduler
# crons single-firing (the scheduler is NOT a separate supervised service).
cd /app
# Invoke via python -m so we don't depend on the uvicorn console script being
# on PATH (only site-packages are copied from the api-build stage, not bin/).
exec python -m uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
